# xiaoya-emby/rootfs/etc/services.d/alist/run
#!/usr/bin/with-contenv bash

cd /opt/alist

exec /opt/alist/alist server --no-prefix